# apps/worker/crontab
0 5 * * * cleanup
0,15,30,45 * * * * fetch-missing-post-records
5,20,35,50 * * * * classify-unclassified-posts

# dev only
#*/5 * * * * classify-unclassified-posts

